# adpki-installer/package-files/usr/lib/systemd/system/adpki-reverb.service
[Unit]
Description=AD-PKI Laravel Reverb WebSocket Server
After=network-online.target postgresql.service php8.4-fpm.service
Wants=network-online.target postgresql.service

[Service]
Type=simple
User=adpki
Group=adpki
SupplementaryGroups=www-data
WorkingDirectory=/opt/adpki/backend
EnvironmentFile=/etc/adpki/backend.env
ExecStart=/usr/bin/php artisan reverb:start --no-interaction
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/adpki/backend/storage /opt/adpki/backend/bootstrap/cache /var/log/adpki

[Install]
WantedBy=multi-user.target
